# deb-package/DEBIAN/prerm
#!/bin/bash -e

systemctl --global disable nvmf.service || true
systemctl stop nvmf.service || true
