# deb-package/DEBIAN/prerm
#!/bin/bash -e

systemctl --global disable nvmf.service
systemctl stop nvmf.service
