# deb-package/DEBIAN/prerm
#!/bin/bash -e

systemctl disable nvmf.service || true
systemctl stop nvmf.service || true
